# deb_package_sources/noble/switch-l4t-configs/etc/profile.d/firefox.sh
export MOZ_DISABLE_RDD_SANDBOX=1
export MOZ_GFX_SPOOF_GL_VERSION="OpenGL ES 3.2 NVIDIA 480.0.0"
