# deb_package_sources/noble/switch-l4t-configs/usr/lib/systemd/user/onboard-kill.service
[Unit]
Description=Onboard kill
PartOf=graphical-session.target
After=indicators-pre.target
StopWhenUnneeded=yes

[Service]
RemainAfterExit=yes
ExecStop=/usr/bin/killall -9 onboard
